# packaging/debian/postinst
#!/bin/sh
# postinst — Simple Photos
#
# Runs after the .deb is unpacked.  Responsibilities:
#   1. Create the `simple-photos` system user and group (idempotent).
#   2. Scaffold /var/lib/simple-photos/{db,storage,models} with correct perms.
#   3. Generate a random JWT secret on first install (preserve on upgrade).
#   4. Hand off to dh_installsystemd / cargo-deb's systemd handler.
#
# The actual photo storage root is finalised in the first-run setup wizard,
# which can move data to any user-readable path.

set -e

USER=simple-photos
GROUP=simple-photos
DATA=/var/lib/simple-photos
CONF=/etc/simple-photos/config.toml

case "$1" in
    configure)
        # 1. Create system user (locked password, no shell).
        if ! getent group "$GROUP" >/dev/null; then
            addgroup --system "$GROUP"
        fi
        if ! getent passwd "$USER" >/dev/null; then
            adduser --system --ingroup "$GROUP" \
                    --home "$DATA" \
                    --shell /usr/sbin/nologin \
                    --gecos "Simple Photos server" \
                    "$USER"
        fi

        # 2. Data directories.
        mkdir -p "$DATA/db" "$DATA/storage" "$DATA/models" "$DATA/certs" "$DATA/downloads"
        chown -R "$USER:$GROUP" "$DATA"
        chmod 750 "$DATA"

        # 2b. Seed the bundled Android APK into the writable downloads dir so
        #     the in-app "Download APK" button works immediately and offline
        #     (issue #4). The package ships a read-only copy under
        #     /usr/share/simple-photos; the server serves it from
        #     {working-dir}/downloads (= $DATA/downloads). The setup service's
        #     fetch-assets.sh later refreshes it from the GitHub release.
        #     Only seed if the downloads copy is missing/empty so we never
        #     clobber a release-fetched (possibly newer) APK on upgrade.
        BUNDLED_APK=/usr/share/simple-photos/simple-photos.apk
        if [ -s "$BUNDLED_APK" ] && [ ! -s "$DATA/downloads/simple-photos.apk" ]; then
            cp "$BUNDLED_APK" "$DATA/downloads/simple-photos.apk"
            chown "$USER:$GROUP" "$DATA/downloads/simple-photos.apk"
            echo "simple-photos: seeded bundled Android APK into $DATA/downloads"
        fi

        # 3. Random JWT secret on first install (only if the placeholder is
        #    still present — we must not clobber an admin-rotated secret on
        #    upgrade).
        if [ -f "$CONF" ] && grep -q 'CHANGE_ME_AT_INSTALL_TIME' "$CONF"; then
            SECRET=$(head -c 32 /dev/urandom | od -An -tx1 | tr -d ' \n')
            # SECRET is pure hex (0-9a-f) so '|' is safe as a sed delimiter
            # and avoids escaping issues with the default '/'.
            sed -i "s|CHANGE_ME_AT_INSTALL_TIME|${SECRET}|" "$CONF"
            chown "$USER:$GROUP" "$CONF"
            chmod 640 "$CONF"
            echo "simple-photos: generated random JWT secret in $CONF"
        fi
        # Always tighten config perms (root:simple-photos 640).
        if [ -f "$CONF" ]; then
            chown root:"$GROUP" "$CONF"
            chmod 640 "$CONF"
        fi

        # 4. NVIDIA / CUDA runtime check (informational — not a hard dep).
        #    The server binary includes the CUDA execution provider and will
        #    use GPU inference automatically when libcudart is present.
        #    If nvidia-smi is found but the CUDA runtime is absent, advise
        #    the operator so they can install it.
        if command -v nvidia-smi >/dev/null 2>&1; then
            if ! ldconfig -p 2>/dev/null | grep -q 'libcudart'; then
                echo "simple-photos: NVIDIA GPU detected but CUDA runtime not found."
                echo "  Install the CUDA runtime for GPU-accelerated AI inference:"
                echo "    sudo apt-get install nvidia-cuda-toolkit"
                echo "  (or follow https://developer.nvidia.com/cuda-downloads for the"
                echo "   NVIDIA-published packages, which are usually more up to date)"
            else
                echo "simple-photos: NVIDIA GPU + CUDA runtime detected — GPU inference enabled."
            fi
        fi

        # 5. First-run asset fetch (AI models + GeoNames + Android APK). A
        #    oneshot service does the heavy network download out-of-band so
        #    dpkg never blocks and offline installs still succeed (it retries
        #    on the next boot until one clean run writes the marker file).
        #    --no-block returns immediately; the fetch runs in the background.
        if [ -d /run/systemd/system ]; then
            systemctl enable simple-photos-setup.service >/dev/null 2>&1 || true
            systemctl start --no-block simple-photos-setup.service >/dev/null 2>&1 || true
            echo "simple-photos: fetching AI models, GeoNames dataset and Android APK in the background"
            echo "  (track it with: journalctl -u simple-photos-setup -f)"
        fi

        # 6. Refresh the desktop / icon caches so the launcher entry and icon
        #    appear immediately instead of after the next login. Best-effort:
        #    these tools may be absent on a headless/server install.
        if command -v update-desktop-database >/dev/null 2>&1; then
            update-desktop-database -q /usr/share/applications >/dev/null 2>&1 || true
        fi
        if command -v gtk-update-icon-cache >/dev/null 2>&1; then
            gtk-update-icon-cache -q -t -f /usr/share/icons/hicolor >/dev/null 2>&1 || true
        fi
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
